# systemd/.config/systemd/user/mbsync.timer
[Unit]
Description=sync IMAP mailbox

[Timer]
#OnCalendar=*-*-* *:0/02
OnBootSec=5min
OnUnitActiveSec=2min

[Install]
WantedBy=default.target
